# roles/dockerinstance/templates/supervisor.conf.j2
[program:{{ item.name }}]
command=/usr/local/bin/dockerrun.sh {{ item.name }} /usr/src/{{ item.name }}/data
directory=/
user=root
